# letsencrypt.service
[Unit]
Description=Let's Encrypt renewal

[Service]
ExecStart=/opt/scripts/certbot_renew.sh
